# lib/capistrano/templates/arduino.service.erb
[Unit]
Description=Arduino HTTP Server for <%= "#{fetch(:application)} (#{fetch(:stage)})" %>
After=network.target

[Service]
Type=simple
<%="User=#{arduino_user(@role)}" if fetch(:arduino_systemctl_user) == :system %>
WorkingDirectory=<%= current_path %>
ExecStart= /bin/bash -lc 'RAILS_ENV=<%= fetch(:rails_env) %> bundle exec rails runner lib/arduino_server.rb'

ExecReload=/bin/kill -TSTP $MAINPID
StandardOutput=append:<%= fetch(:arduino_access_log) %>
StandardError=append:<%= fetch(:arduino_error_log) %>

Restart=always
RestartSec=1

SyslogIdentifier=puma

[Install]
WantedBy=<%=(fetch(:arduino_systemctl_user) == :system) ? "multi-user.target" : "default.target"%>